# leste-config-pixel3a/usr/lib/udev/rules.d/90-bootmac-wifi.rules
ACTION=="add", SUBSYSTEM=="bluetooth", KERNEL=="hci0", \
	  RUN+="/usr/bin/bootmac --bluetooth"
