# debian/leste-config-n900.init
#!/sbin/openrc-run

name="leste-config-n900"

start() {
    echo 512 > /sys/block/mmcblk0/queue/nr_requests
    echo 512 > /sys/block/mmcblk1/queue/nr_requests
}
