# leste-config-common/etc/profile.d/qt-maemo.sh.leste
export QT_STYLE_OVERRIDE=maemo5
export QT_QPA_PLATFORM=maemo

# Input module is not yet ready
#export QT_IM_MODULE=him
