# leste-config-common/etc/profile.d/xdg-runtime.sh.leste
#!/bin/sh

if [ "x$(id)" = "xuser" ]; then
	XDG_RUNTIME_DIR="/run/user/$(id -u)"
	export XDG_RUNTIME_DIR
fi
